# systemd/docker-prune.service
[Unit]
Description=Docker prune (images/containers/volumes)

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -af --volumes
